# supervisord.d/redis-queue.ini
[program:redis-queue]
process_name=%(process_num)02d
command=/usr/bin/php artisan queue:work
directory=/home/wwwroot/default/
numprocs=1
startsecs=1
startretries=10
autostart=true
autorestart=true
user=www
redirect_stderr=true
stdout_logfile=/home/wwwroot/default/storage/supervisor-logs/redis-queue.log